# wally-preview.service
[Unit]
Description=wally-preview — SSRF-hardened Matrix URL preview shim
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/wally-preview.env
ExecStart=/usr/local/bin/wally-preview
Restart=on-failure
RestartSec=2

# Hardening. Note: the SSRF guard in the app is the primary defense; these
# sandbox the process. For belt-and-suspenders egress control, firewall the
# unit's traffic to public IPs only (see README).
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
